# {{cookiecutter.repo_name}}/deployment_configs/{{cookiecutter.repo_name}}_supervisor.conf
[program:{{cookiecutter.repo_name}}]
command=gunicorn manage:app -b localhost:5000
directory=/deploy/app